# roles/kopia/templates/kopia_backup_template.timer.j2
[Unit]
Description={{ description }}

[Timer]
OnCalendar={{ time_string }}
Unit={{ unit_name }}.service

[Install]
WantedBy=timers.target
